# deploy/template/etc/systemd/system/docker.service.d/docker.conf
[Service]
ExecStart=
ExecStartPre=/usr/local/bin/configure_volumes
ExecStart=/usr/bin/dockerd -H fd:// -s overlay2 -g /mnt/var/lib/docker
